# deploy/systemd/plantlab-capture.service.template
[Unit]
Description=PlantLab multi-project capture service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__REPO_PATH__
ExecStart=__PNPM_BIN__ camera:service
Restart=on-failure
RestartSec=5
# Loaded if present; missing files are ignored (leading "-").
EnvironmentFile=-__REPO_PATH__/.env.local
EnvironmentFile=-%h/.config/plantlab/capture.env

[Install]
WantedBy=default.target
